# etc/profile.d/xdg_cache_home.sh
#!/bin/bash

export XDG_CACHE_HOME="/var/cache/users/${USER}/.cache"
